# profile.d/20-vars
#!/bin/bash

# default editor
export EDITOR=vim